# implementations/obsidian-agy/.googooli/scripts/googooli.service
[Unit]
Description=Googooli Telegram Gateway
After=network.target

[Service]
Type=simple
WorkingDirectory=/path/to/your/obsidian-vault
ExecStart=/usr/bin/bash /path/to/your/obsidian-vault/.googooli/scripts/run_gateway.sh
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
